# rpmfiles/sinagios.postinstall
# Have to do this since FPM doesn't know how to set ownership on files
/bin/chown sinagios.sinagios /var/log/sinagios

# Set up init script
/sbin/chkconfig --add sinagios
/sbin/chkconfig sinagios on
